# openwrt.service
[Unit]
Description="OpenWRT"

[Service]
Type=simple
ExecStart=

[Install]
WantedBy=network-online.target