# openwrt.service
[Unit]
Description="OpenWRT"
After=network.target

[Service]
Type=simple
ExecStart=

[Install]
WantedBy=network-online.target